# configs/owncloud.cron
*/5  *  *  *  * php -f /usr/share/nginx/html/owncloud/cron.php > /dev/null 2>&1
